# deploy/huggingface-fullstack-space/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:backend]
command=/bin/bash -lc "PORT=${BACKEND_PORT:-3001} node server.js"
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:frontend]
command=/bin/bash -lc "PORT=${FRONTEND_PORT:-3000} NEXT_PUBLIC_API_URL=${NEXT_PUBLIC_API_URL:-/api} NEXT_PUBLIC_HF_SPACE_URL=${NEXT_PUBLIC_HF_SPACE_URL:-/chatbot} HF_SPACE_URL=${HF_SPACE_URL:-http://127.0.0.1:7861} HF_CHAT_API_NAME=${HF_CHAT_API_NAME:-/chat} npm run start"
directory=/app/turbobujias-web
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:chatbot]
command=/bin/bash -lc "PORT=${CHATBOT_PORT:-7861} GRADIO_SERVER_NAME=${GRADIO_SERVER_NAME:-0.0.0.0} GRADIO_MOUNT_PATH=${GRADIO_MOUNT_PATH:-/chatbot} python app.py"
directory=/app/turbobujias-ai
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
